# cloudflare-dns-change.service
[Unit]
Description=Cloudflare DNS

[Service]
ExecStart=/bin/cloudy

[Install]
WantedBy=multi-user.target